# files/etc/cron.d/update-hosts
# Run update-hosts script
*/5 * * * * root /usr/local/bin/update-hosts.sh > /dev/null
